# nvidia-kmod/nvidia-device-setup.service
[Unit]
Description=Create NVIDIA device nodes
DefaultDependencies=no
After=systemd-modules-load.service systemd-udev-settle.service
Before=multi-user.target
# Only run if nvidia module is loaded
ConditionPathExists=/sys/module/nvidia

[Service]
Type=oneshot
# Wait briefly for driver initialization
ExecStartPre=/usr/bin/sleep 2
# Create primary NVIDIA device nodes (nvidiactl, nvidia0, nvidia-modeset)
ExecStart=/usr/bin/nvidia-modprobe -c 0
# Create UVM device nodes
ExecStart=/usr/bin/nvidia-modprobe -u
# Verify device nodes were created
ExecStartPost=/usr/bin/sh -c 'test -e /dev/nvidia0 && test -e /dev/nvidiactl || exit 1'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
